# supervisor_couchdb.conf
[program:couchdb]
command=/usr/bin/couchdb
autostart=true
autorestart=true
stderr_logfile=/var/log/couchdb/supervisor_error.log
stdout_logfile=/var/log/couchdb/supervisor.log